# etc/profile.d/inform-remote-connection.sh
../zsh/zprofile.d/inform-remote-connection.sh